# fiscalize-assets/script/crontab
# EXECUTA CARGA DAS COTAS DIARIAMENTE 1:00AM ##################################################################################
00 01 * * * /root/fiscalize/carga.sh
